# cmd/chowkidar/chowkidar.service
[Unit]
Description=Internet Chowkidar - Monitor Internet Censorship
Documentation=https://github.com/gnulinuxindia/internet-chowkidar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/chowkidar run
Restart=on-failure
RestartSec=10s

# Run as dynamic user for security
DynamicUser=yes
StateDirectory=chowkidar
ConfigurationDirectory=chowkidar

# Override config paths for dynamic user
Environment="HOME=%S/chowkidar"

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
PrivateMounts=yes

[Install]
WantedBy=multi-user.target
